# deploy/kladovaya-backup.service
[Unit]
Description=Backup Kladovaya PostgreSQL to S3
After=postgresql.service

[Service]
Type=oneshot
User=kladovaya
Group=kladovaya
EnvironmentFile=/opt/kladovaya/.env
ExecStart=/opt/kladovaya/backup-database.sh
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
